# monitoring_setup/atlas-performance.timer
[Unit]
Description=Run Atlas Performance Monitoring every 1 minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
