# monitoring_setup/atlas-performance.service
[Unit]
Description=Atlas Performance Monitoring
After=network.target

[Service]
Type=oneshot
User=atlas
WorkingDirectory=/opt/atlas
ExecStart=/opt/atlas/venv/bin/python monitoring_setup/performance_monitor.py
