# monitoring_setup/atlas-health-check.service
[Unit]
Description=Atlas Health Check
After=network.target

[Service]
Type=oneshot
User=atlas
WorkingDirectory=/opt/atlas
ExecStart=/opt/atlas/venv/bin/python monitoring_setup/health_check.py
